# services/platform-os-detection.service.in
[Unit]
Description=@PROJECT_NAME@ - os detection service
DefaultDependencies=no
After=systemd-modules-load.service
After=systemd-remount-fs.service
ConditionPathIsReadWrite=/etc/

[Service]
Type=oneshot
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/pr-detect --step=os --output-file=@DETECTION_OUTPUT_DIR@/os --json-file=@DETECTION_OUTPUT_DIR@/os.json

[Install]
WantedBy=systemd-modules-load.service

